# deploy/systemd/youqinggong-spring.service
[Unit]
Description=YouQingGong Spring Boot Service
After=network.target mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/youqinggong/backend
ExecStart=/usr/bin/java -jar /srv/youqinggong/backend/app.jar --spring.config.additional-location=/etc/youqinggong/application-prod.yml
SuccessExitStatus=143
Restart=always
RestartSec=5
StandardOutput=append:/var/log/youqinggong/spring-app.log
StandardError=append:/var/log/youqinggong/spring-app.log

[Install]
WantedBy=multi-user.target
